# upstart/l3overlay.conf.in
# l3overlay - IPsec overlay network manager

description "IPsec overlay network manager"
author "Callum Dickinson <callum.dickinson@catalyst.net.nz>"

start on runlevel [2345]
stop on runlevel [016]

respawn
kill timeout 30

exec __L3OVERLAYD__
